# alterion-panel.service
[Unit]
Description=Alterion Panel - Server Management Dashboard
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/alterion-panel/docker
Environment="COMPOSE_PROJECT_NAME=alterion_panel"

# Start Docker containers
ExecStart=/usr/bin/docker-compose up -d

# Stop Docker containers
ExecStop=/usr/bin/docker-compose down

# Reload configuration
ExecReload=/usr/bin/docker-compose restart backend

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=5

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
